# docker/obsidian-sync/root/etc/s6-overlay/s6-rc.d/svc-obsidian-sync/run
#!/usr/bin/with-contenv bash

export HOME=/config
exec s6-setuidgid abc ob sync --continuous --path /vault
